# systemd/iii-engine.service
[Unit]
Description=iii engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=iii
Group=iii
EnvironmentFile=-/etc/iii/engine.env
ExecStart=/opt/iii/current/bin/iii-engine --config /etc/iii/engine.yaml
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/var/log/iii /etc/iii
PrivateTmp=true
PrivateDevices=true

CPUQuota=200%
MemoryMax=1G
TasksMax=512

[Install]
WantedBy=multi-user.target
